# extras/wg-resolve/wg-resolve@.timer
[Unit]
Description=Minutely WireGuard Endpoint Resolver

[Timer]
OnCalendar=minutely
AccuracySec=5s

[Install]
WantedBy=timers.target
